# lib/generators/capistrano/generals/templates/puma_init.sh.erb
#! /bin/sh
<%= template_to_s('_head').to_s %>
### BEGIN INIT INFO
# Provides:          puma
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Manage puma
# Description:       Starts and Stops puma
### END INIT INFO
set -e

# Feel free to change any of the following variables for your app:
APP_ROOT=<%= current_path %>
PID=<%= fetch(:puma_pid) %>
AS_USER=<%= fetch(:puma_user) %>
PUMA_ENV="<%= fetch(:puma_env) %>"
CMD="export HOME; true "${HOME:=$(getent passwd "$AS_USER" | cut -d: -f6;)}"; cd $APP_ROOT && $PUMA_ENV <%= bundle_puma('-d -C', fetch(:puma_config)) %>"

set -u

sig () {
  test -s "$PID" && kill -$1 `cat $PID`
}

run () {
  if [ "$(id -un)" = "$AS_USER" ]; then
    eval $1
  else
    su -c "$1" - $AS_USER
  fi
}

case "$1" in
start)
  sig 0 && echo >&2 "Already running" && exit 0
  run "$CMD"
  ;;
stop)
  sig QUIT && exit 0
  echo >&2 "Not running"
  ;;
force-stop)
  sig TERM && exit 0
  echo >&2 "Not running"
  ;;
*)
  echo >&2 "Usage: $0 <start|stop|force-stop>"
  exit 1
  ;;
esac
